# bin/polaris-health
#!/bin/sh
### BEGIN INIT INFO
# Provides:          polaris-health
# Required-Start:    $remote_fs $network $syslog
# Required-Stop:     $remote_fs $network $syslog
# Should-Start:
# Should-Stop:
# Default-Start:     3
# Default-Stop:      0 1 6
# Short-Description: Polaris Health
# Description:       Polaris Health
### END INIT INFO

. /etc/default/polaris

POLARIS_CONTROL_CMD="$POLARIS_INSTALL_PREFIX/bin/polaris-health-control"

doControl()
{
    ret=$($POLARIS_CONTROL_CMD $1)
}

doStatus()
{
    doControl ping ; EC=$?
    if test "$EC" != "0"
    then 
        echo "Not running or error."
        exit 1
    else
        echo "Running."
    fi
}

doStart()
{
    echo "Starting Polaris Health... "

    doControl start ; EC=$?
    if test "$EC" != "0"
    then
        exit 1
    fi

    sleep 2

    doControl ping  ; EC=$?
    if test "$EC" != "0"
    then
        echo "Failed to start."
        exit 1
    else
        echo "Started."
    fi 
}

doStop()
{
    echo "Stopping Polaris Health... "

    doControl stop ; EC=$?
    if test "$EC" != "0"
    then
        exit 1
    else
        sleep 2
        doControl ping  ; EC=$?
        if test "$EC" != "0"
        then
            echo "Stopped."
        else
            echo "Failed to stop."
            exit 1
        fi
    fi
}

case "$1" in
    status)
        doStatus
    ;;

    start)
        doStart
    ;;

    start-debug)
        echo "Starting Polaris Health in debug mode... "

        doControl start-debug
    ;;

    stop)
        doStop
    ;;

    restart)
        doStop
        doStart
    ;;
        
    *)
        echo "polaris-health [status\|start\|start-debug\|stop\|restart]"
    ;;
esac
