# bin/polaris-health
#!/bin/sh
### BEGIN INIT INFO
# Provides:          polaris-health
# Required-Start:    $remote_fs $network $syslog
# Required-Stop:     $remote_fs $network $syslog
# Should-Start:
# Should-Stop:
# Default-Start:     3
# Default-Stop:      0 1 6
# Short-Description: Polaris Health
# Description:       Polaris Health
### END INIT INFO

. /etc/default/polaris

POLARIS_CONTROL_CMD="$POLARIS_INSTALL_PREFIX/bin/polaris-health-control"

doControl()
{
    ret=$($POLARIS_CONTROL_CMD $1)
}

checkConfig()
{
        echo "Checking polaris health configuation.."
        $POLARIS_CONTROL_CMD check-config ; EC=$?
        if test "$EC" != "0"
        then
            echo "Configuration check FAILED"
            exit 1
        fi

        echo "Configuration is validated."
}

case "$1" in
    status)
        doControl ping 2>/dev/null ; EC=$?
        if test "$EC" != "0"
        then
            echo "not running"
            exit 1
        else
            echo "running"
        fi
    ;;

    start)
        echo "Starting Polaris Health... "

        checkConfig

        # ensure polaris is not already running
        doControl ping 2>/dev/null ; EC=$?
        if test "$EC" = "0"
        then
            echo "aborted, already running!"
            exit 1
        fi

        doControl start
        sleep 2
        doControl ping 2>/dev/null ; EC=$?
        if test "$EC" != "0"
        then
            echo "failed to start!"
            exit 1
        else
            echo "done."
        fi
    ;;


    check-config)

        checkConfig

    ;;

    start-debug)
        echo "Starting Polaris Health in debug mode... "
        checkConfig
        $POLARIS_CONTROL_CMD start-debug
    ;;

    stop)
        echo "Stopping Polaris Health... "

        # ensure polaris is running
        doControl ping 2>/dev/null ; EC=$?
        if test "$EC" != "0"
        then
            echo "aborted, not running!"
            exit 1
        fi

        doControl stop
        echo "done."
    ;;

    restart)
        echo "Restarting Polaris Health... "

        checkConfig

        # ensure polaris is running
        doControl ping 2>/dev/null ; EC=$?
        if test "$EC" != "0"
        then
            echo "Polaris is not running, starting.."
            doControl start
            exit 0
        fi

        doControl stop
        sleep 2
        doControl start
        sleep 2
        doControl ping 2>/dev/null ; EC=$?
        if test "$EC" != "0"
        then
            echo "failed to start!"
            exit 1
        else
            echo "done."
        fi
    ;;

    *)
    echo polaris-health [status\|start\|start-debug\|stop\|restart\|check-config]
    ;;

esac
