# bin/polaris-health
#!/bin/sh
### BEGIN INIT INFO
# Provides:          polaris-health
# Required-Start:    $remote_fs $network $syslog
# Required-Stop:     $remote_fs $network $syslog
# Should-Start:
# Should-Stop:
# Default-Start:     3
# Default-Stop:      0 1 6
# Short-Description: Polaris Health
# Description:       Polaris Health
### END INIT INFO

. /etc/default/polaris

POLARIS_CONTROL_CMD="$POLARIS_INSTALL_PREFIX/bin/polaris-health-control"

doControl()
{
    ret=$($POLARIS_CONTROL_CMD $1 2>/dev/null)
}

case "$1" in
    status)
        doControl ping ; EC=$?
        if test "$EC" != "0"
        then 
            echo "not running"
            exit 1
        else
            echo "running"
        fi
    ;;

    start)
        echo -n "Starting Polaris Health... "

        # ensure polaris is not already running
        doControl ping ; EC=$?
        if test "$EC" = "0"
        then
            echo "aborted, already running!"
            exit 1
        fi

        doControl start        
        sleep 2
        doControl ping ; EC=$?
        if test "$EC" != "0"
        then
            echo "failed to start!"
            exit 1
        else
            echo "done."
        fi
    ;;

    start-debug)
        echo "Starting Polaris Health in debug mode... "
        $POLARIS_CONTROL_CMD start-debug
    ;;

    stop)
        echo -n "Stopping Polaris Health... "
        
        # ensure polaris is running
        doControl ping ; EC=$?
        if test "$EC" != "0"
        then
            echo "aborted, not running!"
            exit 1
        fi

        doControl stop
        echo "done."
    ;;

    restart)
        echo -n "Restarting Polaris Health... "

        # ensure polaris is running
        doControl ping ; EC=$?
        if test "$EC" != "0"
        then
            echo "aborted, not running!"
            exit 1
        fi

        doControl stop  
        sleep 2
        doControl start
        sleep 2
        doControl ping ; EC=$?
        if test "$EC" != "0"
        then
            echo "failed to start!"
            exit 1
        else
            echo "done."
        fi
    ;;

    *)
    echo polaris-health [status\|start\|start-debug\|stop\|restart]
    ;;

esac

